# images/flask_tester/conf/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --http :80 --wsgi-file /src/wsgi.py 
